# private_dot_config/systemd/user/backup.service
[Unit]
Description=Does a backup

[Service]
ExecStart=nu %h/.local/share/scripts/backup.nu

